# circadian-ramp.service
[Unit]
Description=Leviton Circadian Brightness Ramp
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/circadian-ramp/venv/bin/python /opt/circadian-ramp/main.py service
ExecReload=/bin/kill -USR1 $MAINPID
EnvironmentFile=/etc/circadian-ramp.env
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
